# branches/build-livecd-hardy/config/etc/pudata/schindler.list
/home/
/root/
/etc/network
/etc/ppp
/etc/resolv.conf
/etc/X11/xinit
/etc/X11/xorg.conf
/etc/cups
/etc/init.d/
/etc/pudata/
